# 100-usb.rules
ACTION=="add", SUBSYSTEM=="usb", DRIVER=="usb", RUN+="/opt/midi_connect.sh"
